# infra/systemd/cert-renew.service
[Unit]
Description=Renew Tailscale TLS certificate and reload code-server

[Service]
Type=oneshot
ExecStart=%h/Desktop/fuwari/scripts/renew-cert.sh
