# infra/systemd/vault-publish.service
[Unit]
Description=Vault publish pipeline (snapshot, gitleaks, convert, push)

[Service]
Type=oneshot
ExecStart=%h/Desktop/fuwari/scripts/vault-publish/publish.sh
